# roles/node-exporter/templates/exporter.service.j2
[Unit]
Description=NodeExporter

[Service]
TimeoutStartSec=0
User={{ exporter_user }}
ExecStart={{ exporter_bin_path }} {{ exporter_options }}

[Install]
WantedBy=multi-user.target
